# podcast-bot.service
[Unit]
Description=Podcast English Learning Telegram Bot
After=network.target

[Service]
Type=simple
User=paul
WorkingDirectory=/home/paul/learning-english-chatbot
EnvironmentFile=/home/paul/learning-english-chatbot/.env
ExecStart=/usr/bin/python3 -m bot.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
